# infra-host/roles/ansible-pull-timer/templates/ansible-pull.service.j2
[Unit]
Description=Run AI appliance host converge

[Service]
Type=oneshot
ExecStart={{ ai_appliance_converge_path }}
