# magic-host/roles/ansible-pull-timer/templates/ansible-pull.timer.j2
[Unit]
Description=Run AI appliance host converge every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
